# bmw_m3_daq/scripts/systemd/daq.service
[Unit]
Description=BMW M3 Data Acquisition System
After=network.target gpsd.service bluetooth.service
Wants=gpsd.service bluetooth.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/bmw_m3_daq
ExecStart=/home/pi/bmw_m3_daq/venv/bin/python /home/pi/bmw_m3_daq/src/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
